# src/main/resources/static/cloudwatch.path
[Unit]
Description=Watch for the presence of /opt/amazon-cloudwatch-agent.json

[Path]
PathExists=/opt/amazon-cloudwatch-agent.json

[Install]
WantedBy=cloud-init.target