# src/main/resources/static/csye6225.path
[Unit]
Description=Watch for the presence of /opt/application.properties

[Path]
PathExists=/opt/application.properties

[Install]
WantedBy=cloud-init.target
